# ansible/roles/users/templates/etc/profile.d/init.sh.j2
export HISTTIMEFORMAT="%Y-%m-%d %T "
HISTSIZE=10000
HISTFILESIZE=20000
